# deployment/voxbox@.service
[Unit]
Description=VoxBox - Video to Obsidian Knowledge Pipeline
Documentation=https://github.com/stefanahman/voxbox
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/opt/voxbox

# Start the service
ExecStart=/usr/local/bin/docker-compose -f docker-compose.prod.yml up

# Stop the service
ExecStop=/usr/local/bin/docker-compose -f docker-compose.prod.yml down

# Restart policy
Restart=always
RestartSec=10

# Resource limits
LimitNOFILE=65536

# Security
PrivateTmp=yes
NoNewPrivileges=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voxbox

[Install]
WantedBy=multi-user.target

